# files/homeschule/.config/autostart/firstboot.desktop
[Desktop Entry]
Type=Application
Exec="/opt/setup/setup.sh"
Hidden=false
NoDisplay=false
X-GNOME-Autostart-enabled=true
Name=Startup Script